# config/supervisor.conf
[program:bbs]
command=/usr/local/bin/gunicorn wsgi -c config/gunicorn.py
directory=/var/www/bbs
autostart=true
autorestart=true
